# roles/ipnat/templates/setup-ipnat.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='site/setup/ipnat' type='service' version='0'>
    <create_default_instance enabled='true'/>
    <single_instance/>
    <dependency name='filesystem' grouping='require_all' restart_on='restart' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='setup-init' grouping='require_all' restart_on='restart' type='service'>
      <service_fmri value='svc:/site/setup/init'/>
    </dependency>
    <exec_method name='start' type='method' exec='{{ binary_install_dir }}/setup-ipnat.sh %m' timeout_seconds='10'/>
    <exec_method name='stop' type='method' exec='{{ binary_install_dir }}/setup-ipnat.sh %m' timeout_seconds='10'/>
    <exec_method name='refresh' type='method' exec='{{ binary_install_dir }}/setup-ipnat.sh %m' timeout_seconds='10'/>
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='transient'/>
    </property_group>
    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>IPNAT custom configuration</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
